# data/ubuntu-touch-session/anbox-session.conf
description "Anbox session manager"

start on started unity8 and file FILE=/home/phablet/anbox-data/.enable
stop on stopping unity8

kill timeout 120
respawn
respawn limit unlimited

pre-start script
  initctl set-env ANBOX_NO_DESKTOP_SUBDIR=1
end script

exec anbox session-manager --desktop_file_hint=/usr/share/applications/anbox.desktop --mir-rootless
