# data/ubuntu-touch-session/anbox-container.conf
description "Anbox container manager"

start on started lightdm and file FILE=/home/phablet/anbox-data/.enable

kill timeout 120
respawn
respawn limit unlimited

post-stop script
  anbox-bridge stop
end script

pre-start script
  restart cgroupfs-mount
  chmod 666 /dev/binder1 || true
  chmod 666 /dev/anbox-binder || true
  anbox-bridge start
end script

exec anbox container-manager --data-path=/home/phablet/anbox-data --privileged
