# data/anbox.service.in
[Unit]
Description=Anbox session manager

[Service]
Type=simple
ExecStart=@SNAP_MOUNT_DIR@/anbox session-manager
Restart=always
StartLimitIntervalSec=3
StartLimitBurst=1

[Install]
WantedBy=default.target
